# deploy/d-brain-process.service
[Unit]
Description=d-brain Daily Processing

[Service]
Type=oneshot
User=shima
WorkingDirectory=/home/shima/projects/agent-second-brain
ExecStart=/home/shima/projects/agent-second-brain/scripts/process.sh
Environment=PYTHONUNBUFFERED=1
