# deploy/d-brain-weekly.service
[Unit]
Description=d-brain Weekly Digest
After=network.target

[Service]
Type=oneshot
User=shima
WorkingDirectory=/home/shima/projects/agent-second-brain
ExecStart=/home/shima/.local/bin/uv run python scripts/weekly.py
EnvironmentFile=/home/shima/projects/agent-second-brain/.env
StandardOutput=journal
StandardError=journal
